# job_portal_server
#!/bin/sh
#
# chkconfig: 2345 80 30
# description: Job Portal Server

rails s -d
